# configs/openshift_install.service
[Unit]
Description=install_openshift

[Service]
Type=simple
ExecStart=/root/install_script.sh
